# matrixio-malos.service
[Unit]
Description=Matrix Abstraction Layer for OS
After=matrixio-devices-firmware.service

[Service]
Type=simple
ExecStart=/usr/bin/malos

[Install]
WantedBy=multi-user.target
