# debian/postinst
#!/bin/bash -e

echo "Enabling MALOS service"
systemctl enable matrixio-malos

echo "Please restart your Raspberry Pi after installation"
