# debian/prerm
#!/bin/bash -e

echo "Removing MALOS service"
systemctl stop matrixio-malos
systemctl disable matrixio-malos
systemctl daemon-reload
systemctl reset-failed


