# debian/prerm
#!/bin/bash -e

echo "Removing MALOS service"
systemctl stop matrix-creator-malos
systemctl disable matrix-creator-malos
systemctl daemon-reload
systemctl reset-failed


